# packages/neovim/files/etc/profile.d/zz-nvim-default-editor.sh
# Set vim as default editor
# We use `zz-` as prefix to make sure this is processed last in order to
# override any previously set defaults.

export EDITOR=nvim
